# data/xdg-desktop-portal-swaysettings.service.in
[Unit]
Description=SwaySettings Portal

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.swaysettings
ExecStart=@LIBEXECDIR@/xdg-desktop-portal-swaysettings
